# _cms/_dev.to_hosting-wordpress-over-https-with-docker/setup/revproxy/crontab
# Renew Let's Encrypt SSL Certificates that have < 30 days to go,
# in the morning (UTC time).
0 11 * * * /usr/bin/certbot renew --quiet
